# terraform/files/template/helmcharts/cluster-class/files/etcd-defrag.timer
[Unit] 
Description=Run etcd-defrag.service every day
After=network.target
[Timer]
OnCalendar=*-*-* 02:30:00
[Install]
WantedBy=timers.target